# etc/systemd/system/greetd.service.d/nvidia-drm-wait.conf
[Unit]
# Ensure NVIDIA DRM device is initialized before greetd starts niri (the greeter).
# Without this, niri may start before nvidia-drm registers /dev/dri/card*,
# causing niri to fall back to the AMD GPU (no monitors connected) → black screen.

[Service]
# Poll for the NVIDIA DRM by-path symlink (up to 10s), then proceed regardless.
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15 16 17 18 19 20; do [ -e /dev/dri/by-path/pci-0000:01:00.0-card ] && exit 0; sleep 0.5; done; exit 0'
