# etc/udev/rules.d/99-keyball44.rules
# Keyball44: reset on enumeration to fix bad state at boot/resume.
# Script uses a flag file to prevent recursive triggering from the reset itself.
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="5957", ATTR{idProduct}=="0400", \
    RUN+="/etc/udev/scripts/keyball44-reset.sh $env{DEVPATH}"
