# etc/udev/rules.d/99-network-tuning.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="enp5s0", RUN+="/usr/bin/ethtool -K $name lro on", RUN+="/usr/bin/ethtool -G $name rx 8184"
